# etc/alois/prisma/cron.daily
#!/bin/sh
/usr/sbin/prisma-admindb --clean `date --date="last month" +%d.%m.%Y`
