# etc/cron.d/alois-sink
*/15 * * * *	root	/usr/bin/prisma-cleanup-raws syslogd_raws > /dev/null
