# etc/cron.d/alois-prisma
15 1 * * *      root    /usr/bin/prisma-compress-archives 2>&1 1>&2 | grep -v RAILS 2>/dev/null

