# debian/alois-sink.conffiles
/etc/cron.d/alois-sink
/etc/alois/munin-sink.conf
/etc/init.d/prisma-mysqlpiped
